# resources/supervisor/conf.d/nginx.conf
; This file responsible to start NGINX on the foreground

[program:nginx]
command=nginx
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
